# remailer.conf
description "savvy_mailer container"
author "Paul Scarrone"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a savvy_mailer
end script
